# setup/mac-mini-ingest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kp.claude-ingest</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/opt/knowledge-portal/current/setup/claude-ingest.sh</string>
  </array>

  <key>StartInterval</key>
  <integer>600</integer>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/graywzc/Library/Logs/kp-claude-ingest.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/graywzc/Library/Logs/kp-claude-ingest.err</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/graywzc</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
  </dict>
</dict>
</plist>
